# speedtest.service
[Unit]
Description=AKA Python Speedtest service

[Service]
Type=simple
WorkingDirectory=/home/root/speedtest/
ExecStart=/usr/bin/python speedtest-cli-mod.py &
KillMode=process
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
