# speedtest.service
[Unit]
Description=Start Python Speedtest service

[Service]
Type=forking
WorkingDirectory=/home/root/speedtest_project/speedtest/
ExecStart=/usr/bin/python /home/root/speedtest_project/speedtest/speedtest-cli-mod.py
KillMode=process

[Install]
WantedBy=multi-user.target
